# debian/iio-sensor-proxy.init
#!/sbin/openrc-run

depend() {
        need eudev dbus
        after kmod
}

start() {
        ebegin "Starting iio-sensor-deamon"
        start-stop-daemon --start --background --exec /usr/sbin/iio-sensor-proxy
        eend $?
}

stop() {
        ebegin "Stopping iio-sensor-deamon"
        start-stop-daemon --stop --quiet --exec /usr/sbin/iio-sensor-proxy
        eend $?
}
